# scripts/playwright/chromium-cdp.service
[Unit]
Description=Chromium CDP Browser (headless, remote debugging)
After=network.target

[Service]
User=drosner
ExecStart=/usr/bin/chromium \
  --remote-debugging-port=9222 \
  --remote-debugging-address=127.0.0.1 \
  --no-first-run \
  --no-default-browser-check \
  --headless=new \
  --disable-gpu \
  --single-process
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
